# ops/calcpro-gold-intelligence.service
[Unit]
Description=CalcPro HyperGold Intelligence
After=network-online.target calcpro-quote-relay.service
Wants=network-online.target calcpro-quote-relay.service

[Service]
Type=simple
User=fundingpips-deploy
Group=fundingpips-deploy
WorkingDirectory=/opt/fundingpips-calc-intelligence
ExecStart=/usr/bin/node --disable-warning=ExperimentalWarning /opt/fundingpips-calc-intelligence/intelligence/index.js
Environment=NODE_ENV=production
Environment=INTELLIGENCE_HOST=127.0.0.1
Environment=INTELLIGENCE_PORT=8788
Environment=INTELLIGENCE_DB_PATH=/var/lib/calcpro-intelligence/hypergold.sqlite
Environment=INTELLIGENCE_SEED_PATH=/var/lib/calcpro-intelligence/seed-wallets.json
Environment=INTELLIGENCE_MAX_REQUESTS_PER_MINUTE=600
Environment=BYBIT_RELAY_URL=http://127.0.0.1:8787/api/quotes
Restart=always
RestartSec=5
TimeoutStopSec=20
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/var/lib/calcpro-intelligence
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallArchitectures=native
MemoryMax=384M
TasksMax=64
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
